# deploy/ec2/lecrev-firecracker-network.service
[Unit]
Description=Lecrev Firecracker Network Bootstrap
After=network-online.target
Wants=network-online.target
Before=lecrev-node-agent.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/lecrev/node-agent.env
Environment=APP_USER=lecrev
ExecStart=/usr/local/bin/lecrev-configure-firecracker-network
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
